# deploy/vigie-immo.service
[Unit]
Description=Vigie-Immo — Analyse immobilière
After=network.target postgresql.service

[Service]
User=gouroug
WorkingDirectory=/var/www/vigie-immo/vigie-immo-backend
EnvironmentFile=/var/www/vigie-immo/vigie-immo-backend/.env
ExecStart=/var/www/vigie-immo/vigie-immo-backend/venv/bin/gunicorn \
    --config /var/www/vigie-immo/deploy/gunicorn.conf.py app:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vigie-immo

[Install]
WantedBy=multi-user.target
